# scripts/career-scan.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  career-scan.xml — Windows Task Scheduler definition for career-ops daily scan
  ============================================================================
  Runs: node scripts/auto-scan.mjs --greenhouse-only
  Schedule: 6:00 AM Central Time (CT = UTC-6 standard / UTC-5 daylight)
            The trigger uses local time; Windows respects the machine's timezone.
  Working directory: C:\Users\mattm\career-ops\

  IMPORT INSTRUCTIONS
  -------------------
  Open an elevated (Admin) Command Prompt and run:

      schtasks /create /xml "C:\Users\mattm\career-ops\scripts\career-scan.xml" /tn "Career-Ops Scan"

  To verify the task was created:
      schtasks /query /tn "Career-Ops Scan" /fo LIST

  To run immediately (for testing):
      schtasks /run /tn "Career-Ops Scan"

  To delete the task:
      schtasks /delete /tn "Career-Ops Scan" /f

  NOTE: The task runs as the current user (interactive). If the machine is
  asleep or off at 6:00 AM, Windows will run the task at next login
  (StartWhenAvailable = true).
  ============================================================================
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>
    <Date>2026-04-07T00:00:00</Date>
    <Author>mattm</Author>
    <Description>Daily Greenhouse portal scan for new data/analytics/AI roles. Appends new matches to data/pipeline.md and data/scan-history.tsv.</Description>
    <URI>\Career-Ops Scan</URI>
  </RegistrationInfo>

  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-04-08T06:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <!-- LogonType 7 = TASK_LOGON_INTERACTIVE_TOKEN (runs as logged-in user, no stored password needed) -->
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <!-- node.exe full path avoids PATH lookup issues in non-interactive Task Scheduler context -->
      <Command>C:\Program Files\nodejs\node.exe</Command>
      <Arguments>scripts\auto-scan.mjs --greenhouse-only</Arguments>
      <WorkingDirectory>C:\Users\mattm\career-ops\</WorkingDirectory>
    </Exec>
  </Actions>

</Task>
